# bin/obot.cron
#!/bin/sh

WORKDIR="/home/bart/.obot"
PID=`cat $WORKDIR/obot.pid`
echo $PID
if test -r $PID
then
  if `kill -CHLD $PID 1>/dev/null 2>&1`
  then
    exit 0
  else
    obot -m irc,rss,shw,opr,mods -o suicide -l debug -b 1>/dev/null 2>&1
  fi
fi
